# etc/systemd/system/cronie.service.d/mailoff.conf
[Service]
ExecStart=
# do NOT send mail from cron
# -n: foreground
# -s: syslog
# -m: mail
ExecStart=/usr/bin/crond -n -s -m off
