# etc/systemd/system/vnstat.service.d/setconfig.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/vnstatd --config /etc/vnstat.conf --nodaemon