# etc/systemd/system/nix-daemon.service.d/quota.conf
[Service]
CPUQuota=700%
MemoryHigh=5G
MemoryMax=6G
