# etc/systemd/system/nix-daemon.service.d/quota.conf
[Service]
CPUQuota=600%
MemoryHigh=3G
MemoryMax=5G
