# etc/systemd/system/systemd-suspend.service.d/10-dont-freeze-user-session.conf
## fixes systemd#{33083,33626}
[Service]
Environment=SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=false
